# systemd/airpi.service
[Unit]
Description=AirPI Inference Server
Documentation=https://github.com/CyberG3niusIT/AirPI
After=network.target
StartLimitIntervalSec=60s
StartLimitBurst=3

[Service]
Type=simple
User=alex
Group=alex
WorkingDirectory=/home/alex/AirPI

ExecStart=/home/alex/AirPI/.venv/bin/uvicorn server:app \
    --host 127.0.0.1 \
    --port 11435 \
    --workers 1 \
    --log-level info \
    --no-access-log

Environment=AIRPI_MODELS_DIR=/data/models
Environment=AIRPI_DEFAULT_MODEL=qwen2.5-coder-1.5b-instruct-q4_k_m.gguf
Environment=AIRPI_LARGE_MODEL=gemma-3-4b-it-Q4_K_M.gguf
Environment=AIRPI_N_THREADS=4
Environment=AIRPI_N_CTX_SMALL=4096
Environment=AIRPI_N_CTX_LARGE=2048
Environment=AIRPI_MMAP=true
Environment=AIRPI_MLOCK=false
Environment=AIRPI_MAX_QUEUE=10
Environment=AIRPI_KEEP_ALIVE_TIMEOUT=900
Environment=AIRPI_SPECULATIVE=true
Environment=AIRPI_SPECULATIVE_DRAFT_MODEL=qwen2.5-coder-1.5b-instruct-q4_k_m.gguf

Restart=on-failure
RestartSec=5s
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
